# deploy/systemd/pidash-speedtest.timer
[Unit]
Description=Run the pi-dashboard speedtest every 30 minutes

[Timer]
OnBootSec=3min
OnUnitActiveSec=30min
Persistent=true

[Install]
WantedBy=timers.target
